# inittab

# start login shells
tty1::askfirst:-/bin/sh
tty2::respawn:/bin/getty 115200 tty2
tty3::respawn:/bin/getty 115200 tty3
ttyS0::respawn:/bin/getty 115200 ttyS0

# set hostname
::sysinit:/bin/hostname -F /etc/hostname

# set system time assuming hardware clock is in utc (for iptables -m time --kerneltz)
::sysinit:/bin/hwclock -s -u

# start syslog
::sysinit:/bin/syslogd -b 99

# enable forwarding between network interfaces
::sysinit:/bin/echo 1 > /proc/sys/net/ipv4/ip_forward

# configure routing rules
::sysinit:/etc/iptables-rules.sh

# configure lan ip address and start dhcp server on lan
::sysinit:/etc/lan-ip-config.sh
::sysinit:/bin/udhcpd -S /etc/udhcpd.conf

# clean up any left over process identifiers
::sysinit:/bin/rm -f /var/run/*.pid

# enable network interfaces, and keep them up
::sysinit:/bin/ifplugd -p -i wan0 -I -r /etc/wan-updown-event.sh
::sysinit:/bin/ifplugd -p -i lan0 -I -r /etc/lan-updown-event.sh

# start dhcp client on wan
::respawn:/bin/udhcpc -f -i wan0 -p /var/run/udhcpc.pid -s /etc/wan-dhcp-event.sh > /dev/null 2>&1

# start inetd - telnet, ftp and web servers
::sysinit:/bin/inetd /etc/inetd.conf

# update external ip with ddns service
# ::respawn:/etc/ddns-updater.sh

# keep time using a public time server
::once:/bin/ntpd -p time.google.com

